# flask_service.service
[Unit]
Description=AgriShield Flask ML API
After=network.target mysql.service

[Service]
Type=simple
User=www-data
WorkingDirectory=/path/to/Proto1/AgriShield_ML_Flask
Environment="PATH=/usr/bin:/usr/local/bin"
ExecStart=/usr/bin/python3 /path/to/Proto1/AgriShield_ML_Flask/app.py
Restart=always
RestartSec=10
StandardOutput=append:/path/to/Proto1/AgriShield_ML_Flask/flask_service.log
StandardError=append:/path/to/Proto1/AgriShield_ML_Flask/flask_service.log

[Install]
WantedBy=multi-user.target

